# pi_bird_cam/services/bird-camera@.service
[Unit]
Description=Bird Camera Trigger Service
After=pigpiod.service
Requires=pigpiod.service

[Service]
Type=simple
# Use the user that runs the service
User=%i
Group=%i
# Use explicit path expansion
WorkingDirectory=/home/%i/backyard_bird_cam
ExecStart=/usr/bin/python3 /home/%i/backyard_bird_cam/scripts/simple_pir_trigger.py
# Add environment variables for debugging
Environment=PYTHONUNBUFFERED=1
Environment=SERVICE_DEBUG=1
# Ensure proper permissions
UMask=0022
# Restart configuration
Restart=always
RestartSec=30
# Output directly to journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 