# pi_bird_cam/services/bird-camera@.service
[Unit]
Description=Bird Camera Trigger Service
After=pigpiod.service
Wants=pigpiod.service

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/backyard_bird_cam
ExecStart=/usr/bin/python3 /home/%i/backyard_bird_cam/scripts/simple_pir_trigger.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bird-camera
Environment=PYTHONUNBUFFERED=1
Environment=DEBUG=1

# Security settings
NoNewPrivileges=no
ProtectHome=no
ProtectSystem=no
ReadOnlyDirectories=/
ReadWriteDirectories=/home/%i/backyard_bird_cam
ReadWriteDirectories=/home/%i/backyard_bird_cam/images

# Hardware access
SupplementaryGroups=gpio video

[Install]
WantedBy=multi-user.target 