# pi_bird_cam/services/bird-camera@.service
[Unit]
Description=Bird Camera Trigger Service
After=pigpiod.service
Wants=pigpiod.service

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/backyard_bird_cam
ExecStart=/usr/bin/python3 /home/%i/backyard_bird_cam/pi_bird_cam/scripts/simple_pir_trigger.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bird-camera
Environment=PYTHONUNBUFFERED=1
Environment=DEBUG=1
Environment=LOG_LEVEL=DEBUG
Environment=LOG_FORMAT=%(asctime)s - %(name)s - %(levelname)s - %(message)s

# Security Settings
NoNewPrivileges=no
ProtectHome=no
ProtectSystem=no
PrivateTmp=no
ReadOnlyDirectories=/
ReadWriteDirectories=/home/%i
ReadWriteDirectories=/home/%i/backyard_bird_cam
ReadWriteDirectories=/home/%i/backyard_bird_cam/images

# Hardware Access
SupplementaryGroups=gpio video

# Additional Security Settings
CapabilityBoundingSet=
AmbientCapabilities=
DynamicUser=no

# Logging Configuration
LogLevelMax=debug
LogExtraFields=uid=%u gid=%g
LogRateLimitIntervalSec=0
LogRateLimitBurst=0

[Install]
WantedBy=multi-user.target 