# scripts/service/bird_camera.service
[Unit]
Description=Bird Camera System
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/bird_cam
ExecStart=/home/pi/bird_cam/scripts/launch_bird_cam.sh
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=birdcam

[Install]
WantedBy=multi-user.target 