# scripts/service/bird_cam.service
[Unit]
Description=Bird Camera Motion Detection Service
After=network.target pigpiod.service

[Service]
User=%USER%
WorkingDirectory=%HOME%/bird_cam
ExecStartPre=/usr/bin/sudo /usr/bin/systemctl start pigpiod
ExecStart=/usr/bin/python3 %HOME%/bird_cam/bird_cam.py
EnvironmentFile=-%HOME%/bird_cam/.env
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 