# scripts/service/pigpiod.service
[Unit]
Description=Daemon required to control GPIO pins via pigpio
After=network.target

[Service]
Type=forking
User=root
ExecStart=/usr/bin/pigpiod
ExecStop=/bin/systemctl kill pigpiod
Restart=always
RestartSec=10
# Give time for the daemon to fully initialize
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target 