# scripts/service/pigpiod.service
[Unit]
Description=PIGPIO Daemon Service for Bird Camera
After=multi-user.target
Before=bird-camera.service

[Service]
Type=forking
ExecStart=/usr/bin/pigpiod
ExecStop=/bin/systemctl kill pigpiod
Restart=always
RestartSec=10
StandardOutput=append:/var/log/pigpiod.log
StandardError=append:/var/log/pigpiod.error.log

[Install]
WantedBy=multi-user.target 